# deploy/systemd/familycam-main-app.service
[Unit]
Description=FamilyCam Main App (MJPEG)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=orangepi
Group=video
WorkingDirectory=/opt/familycam/projects/main-app
EnvironmentFile=-/opt/familycam/projects/main-app/.env
ExecStart=/opt/familycam/projects/main-app/target/release/family-cam-on-demand-home-monitor
Restart=always
RestartSec=5
TimeoutStopSec=15

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectKernelModules=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryMax=256M
LimitNOFILE=1024

[Install]
WantedBy=multi-user.target
