# deploy/systemd/familycam-go-pion.service
[Unit]
Description=FamilyCam Go Pion Beta
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=orangepi
Group=video
WorkingDirectory=/opt/familycam/projects/beta-go-pion
EnvironmentFile=-/opt/familycam/projects/beta-go-pion/.env
ExecStart=/opt/familycam/projects/beta-go-pion/familycam-go-pion-poc
Restart=always
RestartSec=5
TimeoutStopSec=15

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectKernelModules=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryMax=384M
LimitNOFILE=2048

[Install]
WantedBy=multi-user.target
