# deploy/systemd/familycam-webrtc-rs.service
[Unit]
Description=FamilyCam WebRTC Rust Beta
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=orangepi
Group=video
WorkingDirectory=/opt/familycam/projects/beta-webrtc-rs/webrtc-rs-poc
EnvironmentFile=-/opt/familycam/projects/beta-webrtc-rs/webrtc-rs-poc/.env
ExecStart=/opt/familycam/projects/beta-webrtc-rs/webrtc-rs-poc/target/release/familycam-webrtc-rs-poc
Restart=always
RestartSec=5
TimeoutStopSec=15

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectKernelModules=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryMax=384M
LimitNOFILE=2048

[Install]
WantedBy=multi-user.target
